# supervisord.conf
; ─────────────────────────────────────────────────────────
; NeoMind Supervisord Configuration
;
; Manages four processes inside the Docker container:
;   1. neomind-agent   — Main agent process (CLI or Telegram)
;   2. health-monitor  — Heartbeat detection + HTTP health endpoint
;   3. watchdog        — Last-resort hang detection + forced restart
;   4. data-collector  — 24/7 background financial data collection
;
; Process hierarchy:
;   tini (PID 1) → supervisord → {agent, health-monitor, watchdog, data-collector}
;
; Process independence:
;   - agent restart → collector unaffected
;   - collector restart → agent unaffected
;   - SQLite WAL shared: one writer (collector), many readers (agent)
;
; Usage:
;   supervisord -c /etc/supervisor/conf.d/neomind.conf
;
; Inside container:
;   supervisorctl status           # Check all processes
;   supervisorctl restart agent    # Restart the agent
;   supervisorctl tail agent       # View agent logs
; ─────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/data/neomind/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
pidfile=/tmp/supervisord.pid
childlogdir=/data/neomind

; ── Main Agent Process ──────────────────────────────────

[program:neomind-agent]
; NEOMIND_CMD is set by docker-entrypoint.sh based on mode:
;   telegram mode: "python -u -c 'from agent.finance.telegram_bot import ...'"
;   interactive:   "python -u /app/main.py --mode fin"
command=%(ENV_NEOMIND_CMD)s
directory=/app
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=15
redirect_stderr=true
stdout_logfile=/data/neomind/agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app

; ── Health Monitor ──────────────────────────────────────

[program:health-monitor]
command=python -u /app/agent/evolution/health_monitor.py
directory=/app
autorestart=true
startretries=3
startsecs=2
redirect_stderr=true
stdout_logfile=/data/neomind/health.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app

; ── Watchdog ────────────────────────────────────────────

[program:watchdog]
command=python -u /app/agent/evolution/watchdog.py
directory=/app
autorestart=true
startretries=3
startsecs=2
redirect_stderr=true
stdout_logfile=/data/neomind/watchdog.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app

; ── Data Collector (24/7 background) ──────────────────

[program:data-collector]
command=python -u /app/agent/data/collector.py
directory=/app
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/data/neomind/collector.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app

; ── Supervisor RPC (for supervisorctl) ──────────────────

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
